# zigbee2mqtt/rootfs/etc/s6-overlay/s6-rc.d/zigbee2mqtt-migrate-config/up
/etc/s6-overlay/scripts/zigbee2mqtt-migrate-config.sh